# debpack/etc/cron.d/nems-saas
@reboot root /usr/local/bin/nems-saas-probe > /dev/null 2>&1
