# debpack/etc/cron.d/nems-saas
@reboot root /usr/local/bin/nems-saas-connector >> /var/log/nems-saas.log 2>&1
